# systemd/system/snap.amazon-ssm-agent.amazon-ssm-agent.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application amazon-ssm-agent.amazon-ssm-agent
Requires=snap-amazon\x2dssm\x2dagent-1335.mount
Wants=network.target
After=snap-amazon\x2dssm\x2dagent-1335.mount network.target
X-Snappy=yes

[Service]
ExecStart=/usr/bin/snap run amazon-ssm-agent
SyslogIdentifier=amazon-ssm-agent.amazon-ssm-agent
Restart=always
WorkingDirectory=/var/snap/amazon-ssm-agent/1335
TimeoutStopSec=60
Type=simple
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
